# files/uci-defaults/99-custom-init

# Log to the system log and echo if needed
log_say()
{
    SCRIPT_NAME=$(basename "$0")
    echo "${SCRIPT_NAME}: ${1}"
    logger "${SCRIPT_NAME}: ${1}"
}

/etc/init.d/custom enable
/etc/init.d/custom start
log_say "Custom init service enabled and started"
