# files/init.d/custom
#!/bin/sh /etc/rc.common

START=95
STOP=01

start() {
    if [ -f /etc/rc.custom ]; then 
        bash /etc/rc.custom &
    fi
}
